# systemd/vts-diarization.service
[Unit]
Description=VTS Diarization Container
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Diarization config lives in its OWN env file (diar.env), never vts.env: the
# sidecar must stay free of app secrets (VTS_OAUTH_*, etc.). diar.env carries
# DIARIZATION_IMAGE (for the ${...} substitutions in this unit) plus the DIAR_*
# tuning knobs, and is passed whole into the container via --env-file so the
# sidecar reads them. DIAR_CLUSTERING_FA=0.15 there raises VBx acoustic scaling
# off the model's 0.07 so short interjections nested in a longer turn separate
# instead of being smeared into the host speaker (task b07f7491 / vts-55o).
# Same CONT_NAME/VTS_ROLE convention as vts-worker/vts-webapi, so the
# container name and its DNS alias are derived the same way.
Environment=CONT_NAME=vts
Environment=VTS_ROLE=diarization
EnvironmentFile=/opt/vts/config/diar.env
ExecStartPre=/usr/bin/podman pull ${DIARIZATION_IMAGE}
# No prompts/state/data mounts: the sidecar only takes audio over HTTP and
# returns speaker segments. Its weights are baked into the image, so it needs
# no volumes and no network access beyond serving the API.
ExecStart=/usr/bin/podman run --rm \
  --network-alias=${CONT_NAME}-${VTS_ROLE}.dns.podman \
  --name=${CONT_NAME}-${VTS_ROLE} \
  --env-file=/opt/vts/config/diar.env \
  -p 9100:9100 \
  ${DIARIZATION_IMAGE}
ExecStop=/usr/bin/podman stop -t 20 ${CONT_NAME}-${VTS_ROLE}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
